# debian/prerm
#!/bin/sh
# SoundSync prerm — run by dpkg before package removal.
set -e

case "$1" in
    remove|deconfigure)
        # G12: act ONLY on remove|deconfigure — NOT on upgrade/failed-upgrade.
        # On an upgrade the units must keep running (stopping the CoD watcher
        # mid-upgrade lets the adapter drift off its audio CoD), and the new
        # postinst reconciles unit state. Tearing units down here on every
        # upgrade also double-manages the user unit against postinst's
        # re-enable. So upgrade/failed-upgrade fall through to the default no-op.

        # Disable user unit globally (best-effort).
        systemctl --global disable soundsync.service 2>/dev/null || true

        # Stop + disable the system adapter unit (best-effort).
        if systemctl is-enabled soundsync-adapter.service >/dev/null 2>&1; then
            systemctl disable --now soundsync-adapter.service || true
        else
            systemctl stop soundsync-adapter.service 2>/dev/null || true
        fi

        if command -v deb-systemd-helper >/dev/null 2>&1; then
            deb-systemd-helper disable soundsync-adapter.service || true
        fi

        # Stop + disable the Class-of-Device watcher unit (best-effort).
        if systemctl is-enabled soundsync-adapter-cod.service >/dev/null 2>&1; then
            systemctl disable --now soundsync-adapter-cod.service || true
        else
            systemctl stop soundsync-adapter-cod.service 2>/dev/null || true
        fi

        if command -v deb-systemd-helper >/dev/null 2>&1; then
            deb-systemd-helper disable soundsync-adapter-cod.service || true
        fi
        ;;

    *)
        ;;
esac

#DEBHELPER#

exit 0
